# .crontab
@weekly /home/ootput/bin/lbdbfill
@hourly /usr/bin/offlineimap
@hourly /usr/bin/mairix
